# app-learn-api.service
[Unit]
Description=App Learn API Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=applearn
WorkingDirectory=/home/applearn/app-learn-backend
EnvironmentFile=/home/applearn/app-learn-backend/.env.production

# Start command using PM2
ExecStart=/usr/local/bin/pm2 start ecosystem.config.js --no-daemon

# Restart policy
Restart=on-failure
RestartSec=10s

# Graceful shutdown
KillSignal=SIGTERM
KillMode=process
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Process management
StandardOutput=journal
StandardError=journal
SyslogIdentifier=app-learn-api

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/app-learn /home/applearn/app-learn-backend

[Install]
WantedBy=multi-user.target
